# com.heatmap.bot.dashboard.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heatmap.bot.dashboard.local</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- Assuming the path on Pluto is ~/heatmap-of-fascism-bot -->
        <!-- We use a helper script to find the dir dynamically if possible, or hardcode typical structure -->
        <string>/Users/oscar/heatmap-of-fascism-bot/start_dashboard_server.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/oscar/heatmap-of-fascism-bot</string>
    <key>StandardErrorPath</key>
    <string>/Users/oscar/heatmap-of-fascism-bot/logs/dashboard_server_launch.err</string>
    <key>StandardOutPath</key>
    <string>/Users/oscar/heatmap-of-fascism-bot/logs/dashboard_server_launch.log</string>
</dict>
</plist>
